# docker-assets/etc/supervisor.d/appdaemon.ini
[program:appdaemon]
user=thermo
directory=/opt/thermo-center
command=python manage.py appdaemon
