# docker-assets/etc/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off; pid /var/run/nginx.pid;'
