# docker-assets/etc/supervisor.d/uwsgi.ini
[program:uwsgi]
user=thermo
directory=/opt/thermo-center
command=uwsgi --ini uwsgi.ini
autorestart=true
stopsignal=INT
